# install/etc/s6/services/10-clamd/run
#!/usr/bin/with-contenv bash

source /assets/functions/*
PROCESS_NAME="clamav"
check_initialized

### Set Debug Mode
if [ "$DEBUG_MODE"="TRUE" ] || [ "$DEBUG_MODE"="true" ]; then
	DEBUG_MODE="--debug" 
fi

liftoff

print_info "Starting ClamAV"
exec s6-setuidgid clamav clamd $DEBUG_MODE
